# spec-files-extra/trunk/ext-sources/saslauthd.xml
<?xml version='1.0'?>
 <!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
	Created for SFEcyrus-sasl by Ian Johnson 2013-09-05
	ianj@tsundoku.ne.jp
 -->

 <service_bundle type='manifest' name='SFEsaslauthd:default'>

 <service
         name='network/saslauthd'
         type='service'
         version='1'>
         <create_default_instance enabled='false' />
         <single_instance />


         <dependency name='fs'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/filesystem/local' />
         </dependency>

         <dependency name='net'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/loopback' />
         </dependency>

	 <dependent
		name='saslauthd_multi-user'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri
		    value='svc:/milestone/multi-user' />
	 </dependent>

         <exec_method
                 type='method'
                 name='start'
                 exec='/usr/gnu/sbin/saslauthd -a pam'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='stop'
                 exec=':kill'
                 timeout_seconds='-1'>
         </exec_method>



 </service>

 </service_bundle>
